# systemd/tracker-sd-indexing.path
[Unit]
Description=File watcher to run tracker-sd-indexing.
Before=statefs.service

[Path]
PathChanged=/dev/sdcard

[Install]
WantedBy=pre-user-session.target
