# systemd/tracker-sd-indexing.service
[Unit]
Description=Tracker SDCard indexer toggle
Requires=dbus.socket

[Service]
ExecStart=/usr/bin/tracker-sd-indexing.sh


